# system-config/network-monitor-check.timer
[Unit]
Description=Perform network check
Requires=network-monitor-check.service

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target