# system-config/network-monitor-outage-report.timer
[Unit]
Description=Generate outage report for network monitor

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target